# scripts/survey-system.service
[Unit]
Description=Survey System API Server
Documentation=github.com/pvzcs/target-survey-system-backend
After=network.target mysql.service redis.service
Wants=mysql.service redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/survey-system
EnvironmentFile=/opt/survey-system/.env
ExecStart=/opt/survey-system/survey-system
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=survey-system

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/survey-system

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
